# crontab
15 8 * * * /app/update.sh >> /var/log/update.log 2>&1
